# systemd/user/hermes-backup-check.service
# Source template rendered by install.sh into the user's systemd unit dir.
# Secrets stay in the local config file; this unit passes only its path to the script.

[Unit]
Description=Hermes VM restic repository check
Documentation=https://github.com/dannyfranca/hermes-backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=EXAMPLE_REPO_PATH
ExecStart=EXAMPLE_REPO_PATH/scripts/restic-check.sh --config-env EXAMPLE_LOCAL_CONFIG_DIR/hermes-backup.env
